# systems-and-network-adminstration/Week09/updater.service
[Unit]
Description=Service for updating packages

[Service]
ExecStart=/bin/bash apt update

[Install]
WantedBy=multi-user.target
